# dist/init/linux-systemd/bs440.service
[Unit]
Description=Run BS440 bluetooth scale monitor
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
WorkingDirectory=/opt/BS440
ExecStart=/usr/bin/python BS440.py
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
